# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;" -c /app/deploy/nginx.conf
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:uvicorn]
process_name=uvicorn_%(process_num)s
command=sh /app/deploy/run_uvicorn_worker.sh %(process_num)s
numprocs=%(ENV_WORKERS)s
numprocs_start=0
autostart=true
autorestart=true
priority=20
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=TERM
stopasgroup=true
killasgroup=true
